# docker/opcache.ini
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true

[program:queue-worker]
command="php artisan queue:work --queue=default --sleep=3 --tries=3"
autostart=true
autorestart=true
user=www-data

[program:cron]
command=crond -f -l 8
autostart=true
